# daemontools/indieauth/run
#!/bin/sh
# daemontools run script for indieauth.
#
# Runs with the service directory as cwd, so `./env` refers to the envdir
# alongside this script. Uses djb daemontools' own tools -- envdir and
# setuidgid -- not runit's chpst.

exec 2>&1

# Self-heal both directories indieauth needs on every start, not just the
# socket dir -- otherwise a missing /var/lib/indieauth (e.g. before the
# one-time manual setup, or if it's ever wiped) crash-loops instead of
# recovering.
mkdir -p /var/run/indieauth /var/lib/indieauth
chown indieauth:indieauth /var/run/indieauth /var/lib/indieauth

exec envdir ./env setuidgid indieauth /usr/local/bin/indieauth
